# scripts/visioncore.service
[Unit]
Description=VisionCorePro Headless Vision System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teamhd
Group=teamhd
WorkingDirectory=/opt/VisionCorePro
ExecStart=/usr/bin/python3 -m app.main
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
